# home/dot_local/cron.daily/executable_downsort.sh.tmpl
#!/bin/bash

cd {{ .chezmoi.sourceDir }}/../bin/downsort
cargo run --release
